# confs/supervisor-app.conf
[program:app]
command = bash -c "/usr/sbin/nginx"
